# supervisord.conf
[supervisord]
nodaemon=true

[program:vncserver]
user=root
environment=USER="root"
command=/usr/bin/vncserver :1 -geometry 1280x800 -depth 24
autorestart=true

